# common/rootfs_target/etc/profile.d/env.sh
#!/bin/bash

export BARDYLIB=/home/embedded/bardy
export LD_LIBRARY_PATH=${BARDYLIB}
